# roles/gitmonitoring/templates/etc/cron.d/gitcron
# {{ ansible_managed }}
MAILTO={{ report_email_address }}
10 5 * * * root git --git-dir "{{ item.path }}/.git" --work-tree "{{ item.path }}" ls-files -s -t -u -d -m -o --exclude-standard
